# firefox-watchdog.timer
[Unit]
Description=Firefox Kiosk Watchdog Timer
Requires=firefox-watchdog.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target