# shutdown-timer.timer
[Unit]
Description=Timer for automatic system shutdown at 22:00
Requires=shutdown-timer.service

[Timer]
# Run daily at 22:00 (10:00 PM)
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target